# ENskip-0/ENskip-0.10pa/next/rc.local
INSTALLDIR=/usr/local/skip				# skip
if [ -f $INSTALLDIR/rc.skip ]; then			# skip
	$INSTALLDIR/rc.skip start > /dev/console 2>&1	# skip
fi							# skip
